# gbrain-ingest/systemd/gbrain-ingest-drive.service
# systemd USER unit (install: ~/.config/systemd/user/, systemctl --user)
[Unit]
Description=gbrain ingest: recent Google Drive docs -> entity pages
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.hermes/.env
WorkingDirectory=%h/gbrain-ingest
ExecStart=/usr/bin/flock -w 600 /tmp/gbrain-ingest.lock /usr/bin/python3 %h/gbrain-ingest/ingest_drive.py --limit 25
Nice=10
TimeoutStartSec=60min

[Install]
WantedBy=default.target
