# gbrain-ingest/systemd/gbrain-ingest-drive.timer
[Unit]
Description=Run gbrain drive ingestion daily

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=300
Persistent=true
Unit=gbrain-ingest-drive.service

[Install]
WantedBy=timers.target
